# cron.d/update-gs_stats.o.o.sh
#!/bin/bash

pushd /home/dimstar/src/gs-stats.o.o
git pull

